# files/etc/cron.d/checkGlancrUpdates
*/15 * * * * www-data php /var/www/html/cron.php > /dev/null
